# cron.d/safety.cron
@daily /usr/local/bin/check-sudoers.sh | mail -s "Sudoers Diff" tu@mail.com
@hourly chkrootkit > /var/log/chkrootkit_hourly.log
